# ssh-daemon/ssh-tunnel.service
[Unit]
Description=SSH Tunnel Server
After=network.target

[Service]
Type=simple
User=rpal
WorkingDirectory=/var/www/html2/xyemax.test/ssh-daemon
ExecStart=/usr/bin/python3.12 /var/www/html2/xyemax.test/ssh-daemon/server.py \
    --port 2222 \
    --host-key /var/www/html2/xyemax.test/ssh-daemon/host_key
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
